# systemd/media-ingest.service
[Unit]
Description=Media Ingest Server
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/media-ingest -config /etc/media-ingest/config.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
